# dot_config/systemd/user/fetch-wallpaper.timer
[Unit]
Description=Scrape wallpapers from web services

[Timer]
OnUnitActiveSec=10min
OnActiveSec=30min

[Install]
WantedBy=basic.target
