# dot_config/systemd/user/fetch-wallpaper.service
[Unit]
Description=Scrape wallpapers from web services

[Service]
Type=simple
ExecStart=/usr/bin/python .scripts/fetch-wallpaper

